# profile.d/mc-skin.sh
[ ! -t 0 ] && return

[ -z "$MC_SKIN" ] && export MC_SKIN=solarized-light
